# com.maral.jobagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.maral.jobagent.plist
  ─────────────────────────────────────────────────────────
  macOS launchd agent — runs job_agent every hour.

  INSTALLATION:
    1. Edit the paths below to match where you cloned this project.
    2. Copy to LaunchAgents:
         cp com.maral.jobagent.plist ~/Library/LaunchAgents/
    3. Load it:
         launchctl load ~/Library/LaunchAgents/com.maral.jobagent.plist
    4. To unload (stop):
         launchctl unload ~/Library/LaunchAgents/com.maral.jobagent.plist
    5. To run manually right now:
         launchctl start com.maral.jobagent
    6. Check logs:
         tail -f /tmp/job_agent.log
  ─────────────────────────────────────────────────────────
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.maral.jobagent</string>

  <!-- ── EDIT THIS: full path to your Python in the venv ── -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/YOUR_USERNAME/job_agent/venv/bin/python</string>
    <string>/Users/YOUR_USERNAME/job_agent/main.py</string>
  </array>

  <!-- ── EDIT THIS: path to the project folder ── -->
  <key>WorkingDirectory</key>
  <string>/Users/YOUR_USERNAME/job_agent</string>

  <!-- Run every 3600 seconds = every hour -->
  <key>StartInterval</key>
  <integer>3600</integer>

  <!-- Also run immediately when agent is loaded -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Environment variables (API keys go here) -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>ANTHROPIC_API_KEY</key>
    <string>YOUR_ANTHROPIC_API_KEY_HERE</string>
    <key>GMAIL_APP_PASSWORD</key>
    <string>YOUR_GMAIL_APP_PASSWORD_HERE</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/Users/YOUR_USERNAME/job_agent/venv/bin</string>
  </dict>

  <!-- Logs -->
  <key>StandardOutPath</key>
  <string>/tmp/job_agent.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/job_agent_error.log</string>

  <!-- Keep alive on crash -->
  <key>KeepAlive</key>
  <false/>

</dict>
</plist>
